# playbooks/roles/sysctl_ip_forward/templates/systemd/sysctl_ip_forward.timer.j2
[Unit]
Description=Run sysctl-ip-forward.service on schedule

[Timer]
OnBootSec=5min
OnCalendar=*-*-* *:*:00

Unit=sysctl-ip-forward.service

[Install]
WantedBy=timers.target
